# deploy/sports-mcp.service
[Unit]
Description=sports-mcp - MCP server wrapping ESPN's public API for Home Assistant Voice
Documentation=https://github.com/JoeyG1973/sports-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=/opt/sports-mcp

# Optional bind overrides. Defaults are 0.0.0.0:8000.
# Environment=SPORTS_MCP_HOST=0.0.0.0
# Environment=SPORTS_MCP_PORT=8000

# Use the absolute path to uv. Find yours with: command -v uv
# `--no-sync` skips the lockfile-vs-venv check on each start so launches
# are deterministic; re-run `uv sync` manually after pulling updates.
ExecStart=/usr/local/bin/uv run --no-sync sports-mcp

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
